# machines/work-laptop/launchagents/com.timtrailor.work-health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.timtrailor.work-health-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/timtrailor/code/claude-bridge/tools/work_health_check.py</string>
    </array>

    <!-- Hourly probe of work-laptop services. Pattern 3: each probe must
         exercise the actual feature, not just the process. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/work_health_check.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/work_health_check.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/timtrailor/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/timtrailor</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/timtrailor/code/claude-bridge</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
