# machines/work-laptop/launchagents/com.timtrailor.work-credential-rotation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.timtrailor.work-credential-rotation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/timtrailor/code/claude-bridge/tools/work_credential_rotation.py</string>
    </array>

    <!-- Quarterly: 09:00 on Day 1 of Jan/Apr/Jul/Oct.
         launchd fires this entry whenever the next event matches; the
         script gates internally on Month ∈ {1,4,7,10}. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Month</key>
            <integer>1</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>4</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>7</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>10</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/work_credential_rotation.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/work_credential_rotation.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/timtrailor/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/timtrailor</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/timtrailor/code/claude-bridge</string>
</dict>
</plist>
